# provisioning/roles/ansible-ossec/etc/network/if-up.d/authd
#!/bin/sh -e

/var/ossec/bin/ossec-authd >/dev/null 2>&1

exit 0